# etc/cron.d/osh-orphaned-homedir.dist
# On slaves, remove sparse orphaned directories in /home due to deleted accounts/groups
*/5 * * * * root /opt/bastion/bin/cron/osh-orphaned-homedir.sh >/dev/null
